# systemd/obsidian-organize.service.template
[Unit]
Description=Obsidian Daily Organizer

[Service]
Type=oneshot
ExecStart=/bin/bash {{ORGANIZE_SCRIPT}}
Environment=OBSIDIAN_VAULT={{OBSIDIAN_VAULT}}
